# profile.d/runsvdir-user.sh
[ -d "$HOME/service" ] && export SVDIR="$HOME/service/sv"
